# private_dot_config/autostart/private_onedrivegui.desktop
[Desktop Entry]
Categories=Network;Office
Comment=A simple GUI for OneDrive Linux client
Exec=onedrivegui
Icon=OneDriveGUI.png
Name=OneDriveGUI
StartupNotify=true
Type=Application
