# private_dot_config/autostart/rustdesk.desktop

[Desktop Entry]
Type=Application
Exec=rustdesk --tray
NoDisplay=false
        