# private_dot_config/systemd/user/plasma-plasmashell.service.d/override.conf
[Service]
Environment=AMD_VULKAN_ICD=RADV