# private_dot_config/systemd/user/goldwarden.service
[Unit]
Description="Goldwarden daemon"
After=graphical-session.target

[Service]
ExecStart=/usr/bin/goldwarden daemonize

[Install]
WantedBy=graphical-session.target