# private_dot_config/systemd/user/organize.service
[Unit]
Description=Organize files into folders with organize

[Service]
Type=oneshot
ExecStart=organize run
