# private_dot_config/systemd/user/rdiff-bak-tar.service
[Unit]
Description=Backup rdiffbak files into a tar in pCloud

[Service]
Type=oneshot
ExecStart=/home/zany130/Scripts/tarRdiff-Bak.sh
